# udev-rules/10-blockadd.rules
SUBSYSTEM=="block", RUN+="/usr/local/bin/thecus_block_event.sh"

